# gDriveQOAutomation/extension/app/upstart/quickoffice-start.conf
# Copyright 2019 The Chromium Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description "mount Quickoffice squashfs file"
author "chromium-os-dev@chromium.org"

start on started ui

env QO_ROOT="/usr/share/chromeos-assets/quickoffice"

exec mount -o nodev,noexec,nosuid "${QO_ROOT}/quickoffice.squash" \
  "${QO_ROOT}/_platform_specific"
